# systemd/jarvis-self-heal.service
[Unit]
Description=Jarvis self-heal — auto-repair down platforms (guardrailed)
After=network-online.target jarvis-memory.service jarvis-orchestrator.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/jarvis
# Order matters: self-heal.env loads second so its values win any collision
# with secrets.env (this matches the box's real pre-existing unit).
EnvironmentFile=/opt/jarvis/config/secrets.env
EnvironmentFile=/opt/jarvis/config/self-heal.env
ExecStart=/usr/bin/node src/self-heal.js
TimeoutStartSec=120

# Rule 4 safety net (fleet standard) — see jarvis-memory.service for rationale.
# 2026-07-24: this file was reconstructed from the box's REAL unit (found at
# /etc/systemd/system/jarvis-self-heal.service, verified over the tailnet)
# after an earlier same-day guess; only these two limit lines are additions.
MemoryMax=512M
CPUQuota=100%
