# systemd/jarvis-dashboard.service
[Unit]
Description=Jarvis Dashboard — Cyberpunk command center at :9206
After=network.target jarvis-metrics.service jarvis-orchestrator.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/jarvis
ExecStart=/usr/bin/node src/dashboard-server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production
EnvironmentFile=/opt/jarvis/config/secrets.env

[Install]
WantedBy=multi-user.target
