# systemd/jarvis-self-heal.timer
[Unit]
Description=Run Jarvis self-heal controller every 5 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
