# systemd/jarvis-backup.service
[Unit]
Description=Jarvis Memory Backup (WAL-safe SQLite online backup)
After=network.target
Wants=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/jarvis
ExecStart=/opt/jarvis/scripts/backup-memory.sh
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production

# Rule 4 safety net (2026-07-20) — see jarvis-memory.service for rationale.
MemoryMax=512M
CPUQuota=100%
