# systemd/jarvis-fleet-check.timer
[Unit]
Description=Run Jarvis fleet health probe every 10 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
