# systemd/jarvis-agents.service
[Unit]
Description=Jarvis Agent Scheduler — dispatches role agents from config/agents.json (kill switch: AGENTS_MODE)
After=network.target jarvis-memory.service jarvis-orchestrator.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/jarvis
ExecStart=/usr/bin/node src/agent-scheduler.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production
# dry-run: scheduler notifies what it WOULD dispatch but dispatches nothing.
# Craig flips this to AGENTS_MODE=live to activate the org (see CLAUDE.md).
Environment=AGENTS_MODE=dry-run
EnvironmentFile=/opt/jarvis/config/secrets.env

[Install]
WantedBy=multi-user.target
