# systemd/jarvis-slack.service
[Unit]
Description=Jarvis Slack Bridge
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/jarvis
ExecStart=/usr/bin/node src/slack-bridge.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production
EnvironmentFile=/opt/jarvis/config/secrets.env

# Rule 4 safety net (2026-07-20) — see jarvis-memory.service for rationale.
MemoryMax=512M
CPUQuota=100%

[Install]
WantedBy=multi-user.target
